# etc/systemd/madruga-easter.service
[Unit]
Description=Madruga AI Easter — Pipeline orchestrator 24/7
After=network.target

[Service]
Type=notify
ExecStart=/usr/bin/python3 %h/repos/paceautomations/madruga.ai/.specify/scripts/easter.py
WorkingDirectory=%h/repos/paceautomations/madruga.ai
Restart=on-failure
RestartSec=5
WatchdogSec=120
NotifyAccess=all
Environment=PYTHONUNBUFFERED=1
Environment=MADRUGA_MODE=manual

[Install]
WantedBy=default.target
